# scrapy-fang/crontab
0 6 * * * root /scrapy/fang/run_spider.sh
